# config/systemd/qqbot-voice-tts-tailnet.service.example
[Unit]
Description=qqbot laptop-local voice tts tailnet publisher
After=network-online.target tailscaled.service qqbot-voice-tts.service
Wants=network-online.target qqbot-voice-tts.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/kkkzbh/code/qqbot
Environment=QQBOT_VOICE_TTS_TAILNET_ENV_FILE=/home/kkkzbh/code/qqbot/config/voice-tts.tailnet.env
ExecStart=/home/kkkzbh/code/qqbot/scripts/publish-voice-tts-tailnet.sh apply
ExecReload=/home/kkkzbh/code/qqbot/scripts/publish-voice-tts-tailnet.sh apply
ExecStop=/home/kkkzbh/code/qqbot/scripts/publish-voice-tts-tailnet.sh clear

[Install]
WantedBy=default.target
