# config/systemd/qqbot-pmhq.service.example
[Unit]
Description=QQBot PMHQ Service
After=network-online.target
Wants=network-online.target
PartOf=qqbot.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/kkkzbh/code/qqbot
Environment=QQBOT_ENV_FILE=/home/kkkzbh/code/qqbot/.env.local
Environment=CONTAINERS_CONF=/home/kkkzbh/code/qqbot/config/podman/containers.conf
ExecStart=/home/kkkzbh/code/qqbot/scripts/podman-pmhq-service.sh up
ExecStop=/home/kkkzbh/code/qqbot/scripts/podman-pmhq-service.sh stop

[Install]
WantedBy=qqbot.target
