# config/systemd/qqbot-llbot.service.example
[Unit]
Description=QQBot LLBot Service
After=network-online.target qqbot-pmhq.service
Wants=network-online.target qqbot-pmhq.service
PartOf=qqbot.target

[Service]
Type=simple
WorkingDirectory=/home/kkkzbh/code/qqbot
Environment=QQBOT_ENV_FILE=/home/kkkzbh/code/qqbot/.env.local
ExecStart=/home/kkkzbh/code/qqbot/scripts/run-llbot-host.sh
Restart=always
RestartSec=5

[Install]
WantedBy=qqbot.target
