# config/systemd/qqbot-voice-tts.service.example
[Unit]
Description=qqbot laptop-local voice tts gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kkkzbh/code/qqbot
Environment=QQBOT_VOICE_TTS_ENV_FILE=/home/kkkzbh/code/qqbot/config/voice-tts.local.env
ExecStart=/home/kkkzbh/code/qqbot/scripts/run-voice-tts-local.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=default.target
