# config/systemd/qqbot-koishi.service.example
[Unit]
Description=QQBot Koishi Service
After=network-online.target qqbot-llbot.service
Wants=network-online.target qqbot-llbot.service
PartOf=qqbot.target qqbot-llbot.service

[Service]
Type=simple
WorkingDirectory=/home/kkkzbh/code/qqbot
Environment=QQBOT_ENV_FILE=/home/kkkzbh/code/qqbot/.env.local
ExecStart=/home/kkkzbh/.local/bin/pnpm start:local
Restart=always
RestartSec=5

[Install]
WantedBy=qqbot.target
